# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/etc/uci-defaults/uhttpd
#!/bin/sh
#
# Copyright (C) 2014 OpenWrt.org xcloud.cc
#

uci batch <<EOF
set uhttpd.main.script_timeout=600
commit uhttpd
uci commit
EOF

exit 0

