# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/usr/lib/opkg/info/luci-mod-admin-core.list
/etc/init.d/luci_fixtime
/etc/init.d/luci_dhcp_migrate
