# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/usr/lib/opkg/info/mjpg-streamer.list
/usr/lib/output_file.so
/usr/bin/mjpg_streamer
/usr/lib/input_file.so
/usr/lib/output_http.so
/usr/lib/input_uvc.so
/etc/hotplug.d/usb/20-mjpg-streamer
/etc/config/mjpg-streamer
/etc/init.d/mjpg-streamer
