# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/usr/lib/opkg/info/pure-ftpd.list
/usr/bin/pure-statsdecode
/usr/sbin/pure-uploadscript
/usr/sbin/pure-quotacheck
/usr/bin/pure-pwconvert
/etc/config/pure-ftpd
/usr/sbin/pure-mrtginfo
/usr/bin/pure-pw
/usr/sbin/pure-authd
/usr/sbin/pure-ftpwho
/etc/init.d/pure-ftpd
/usr/sbin/ptracetest
/usr/sbin/pure-ftpd
