# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/usr/lib/opkg/info/uhttpd.list
/uhttpd.init
/uhttpd.config
/etc/uci-defaults/uhttpd
/usr/sbin/uhttpd
/etc/init.d/uhttpd
/etc/config/uhttpd
/lib/upgrade/keep.d/uhttpd
