# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/usr/lib/opkg/info/wan-check.list
/etc/init.d/wan-check
/usr/sbin/wancheck
