# firmware report1/test1报告压缩HTML纯静态/firmware/squashfs-root/etc/hotplug.d/usb/20-mjpg-streamer

case "$ACTION" in
	add)
	        # start process
		ls /dev/video0 >/dev/null 2>&1
		if [ $? -ne 1 ]
		then
		/etc/init.d/mjpg-streamer start
		fi
	        ;;
	remove)
	    # stop process
		ls /dev/video0 >/dev/null 2>&1
		if [ $? -ne 0 ]
		then
			/etc/init.d/mjpg-streamer stop
		fi
	        ;;
esac

